# ombot.service
[Unit]
Description=Ombot Relay Service (WebSocket server, Bot side)
After=network.target ombers-communicator.service ombist-openclaw-gateway.service

[Service]
Type=simple
User=ombot
Group=ombot
WorkingDirectory=/opt/ombot/Ombot
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
Environment=PORT=8082
Environment=HEALTH_PORT=9090
Environment=MIDDLEWARE_WS_URL=ws://127.0.0.1:8081/ws
Environment=OPENCLAW_DATA_DIR=/var/lib/ombot
# Optional: merge OpenClaw Gateway (loopback) with Phone via Ombers — set in ombot.env
# OPENCLAW_GATEWAY_BRIDGE=1
# OPENCLAW_BRIDGE_AGENT_ID=default
# OPENCLAW_BRIDGE_CONVERSATION_ID=default
# OPENCLAW_BRIDGE_PARTICIPANT_ID=default
# OPENCLAW_GATEWAY_URL=ws://127.0.0.1:18789
# OPENCLAW_GATEWAY_TOKEN=
EnvironmentFile=-/etc/ombot/ombot.env
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/ombot /opt/ombot/Ombot
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
